# code/agent/service/ga@SENSOR.timer
[Unit]
Description=Growautomation check sensor aht

[Timer]
OnCalendar=*:0/15
Unit=growautomation@SENSOR.service

[Install]
WantedBy=timers.target