# code/agent/service/systemd/ga@DEFAULT.timer
[Unit]
Description=Growautomation check sensor aht

[Timer]
OnCalendar=*:0/15
Unit=ga@SENSOR.service

[Install]
WantedBy=timers.target